# print-queue-manager.service
[Unit]
Description=Print Queue Manager
After=network.target cups.service

[Service]
Type=notify
User=root
WorkingDirectory=/opt/print-queue-manager
Environment="PATH=/opt/print-queue-manager/venv/bin"
EnvironmentFile=/opt/print-queue-manager/.env
ExecStart=/opt/print-queue-manager/venv/bin/gunicorn --bind 0.0.0.0:5000 --workers 2 --timeout 120 app:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
